# services/docker-backup.timer
[Unit]
Description=Weekly Docker Backup

[Timer]
# Run every Monday at 6 AM
OnCalendar=Mon *-*-* 06:00:00

[Install]
WantedBy=timers.target
